# backend/.docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
priority=1
startsecs=0

[program:php-fpm]
command=php-fpm
priority=2
startsecs=0